# files/apt-daily-upgrade.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 2:00
RandomizedDelaySec=1h
